# cyclic-task/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
nodaemon=false
loglevel=info
pidfile=/run/supervisord.pid

[program:task-launcher]
command=/opt/launcher.sh
autostart=true
autorestart=true